# gce-workload-cert-refresh.timer
[Unit]
Description=GCE Workload Certificate refresh timer

[Timer]
OnBootSec=5
OnUnitActiveSec=30m

[Install]
WantedBy=timers.target
