# google-guest-agent.conf
description "GCE Guest Agent"
start on started rsyslog
oom -16

respawn
exec /usr/bin/google_guest_agent
