# google-guest-agent.conf
# Guest agent. Manages network interfaces, user accounts.
start on stopped google-instance-setup
oom -16

respawn
exec /usr/bin/google_guest_agent
