# google-startup-scripts.conf
# Runs a startup script from metadata.
start on started google-guest-agent and startup

exec /usr/bin/google_metadata_script_runner startup
