# google-guest-agent-manager.service
[Unit]
Description=Google Compute Engine Guest Agent Plugin Manager
# Force a hard structural requirement; sshd cannot boot without this service
Requires=sshd-keygen.target
Before=sshd.service sshd-keygen.target
# Group all networking under the official, unified systemd online target
After=network-online.target syslog.service
Wants=network-online.target

[Service]
Type=notify
ExecStart=/usr/bin/google_guest_agent_manager
OOMScoreAdjust=-999
Restart=always
KillMode=process

[Install]
WantedBy=multi-user.target
