# google-guest-agent.conf
description "GCE Guest Agent"
start on stopped rc RUNLEVEL=[2345]
oom -16

respawn
exec /usr/bin/google_guest_agent
